# LowLevel/dispatcher/assets/dispatcher.service
[Unit]
Description=dispatcher of unique number and IO events.
Documentation=https://github.com/rafaelchiafarelli/conboard

[Service]
User=root
# The dispatcher is a long-running daemon (ZMQ + Crow HTTP run loop, joined
# threads), so it must be Type=simple. As Type=oneshot, `systemctl start` waited
# forever for it to exit (installer hang; and a boot-time stall until timeout).
Type=simple
ExecStart=/conboard/LowLevel/dispatcher/build/dispatcher  -x /conboard/LowLevel/dispatcher/assets/config.json
Restart=on-failure


[Install]
WantedBy=multi-user.target
