# LowLevel/assets/100-usb.rules
ACTION=="add", SUBSYSTEM=="usb", DRIVER=="usb", RUN+="/conboard/event_handler.sh"
